# form/web.service
[Unit]
Description=web service
After=network

[Service]
ExecStart=/opt/form

[Install]
WantedBy=multi-user.target
