# f5-version/web.service
[Unit]
Description=web service
After=network

[Service]
ExecStart=/opt/f5web

[Install]
WantedBy=multi-user.target
